# logseq-sync/logseq-sync-push.service
[Unit]
Description=Logseq Notes — Push to remote on shutdown
Before=shutdown.target
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=@@HOME@/.local/share/logseq-sync/logseq-sync-push.sh
TimeoutStopSec=60

[Install]
WantedBy=halt.target poweroff.target reboot.target
