# logseq-sync/logseq-sync-push.timer
[Unit]
Description=Logseq Notes — Scheduled push at 02:00

[Timer]
OnCalendar=*-*-* 02:00:00
Persistent=true
Unit=logseq-sync-push.service

[Install]
WantedBy=timers.target
